# systemd/feedoracle-mcp-risk.service
[Unit]
Description=FeedOracle Stablecoin Risk MCP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/rwa_node/mcp
ExecStart=/usr/bin/python3 /root/rwa_node/mcp/feedoracle_mcp_risk.py
Environment=PORT=5252
Environment=ETHERSCAN_KEY=your_key_here
Environment=FEED_API=http://127.0.0.1:5080
Environment=PEG_API=http://127.0.0.1:5215
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
